# ops/hetzner-update-proxy/drainq-one-mirror.timer
[Unit]
Description=DrainQ ONE — Mirror timer (offset from Suite timer to reduce GitHub API pressure)
Documentation=file:///opt/drainq-mirror-one/DEPLOYMENT.md

[Timer]
# Boot offset 4 min — Suite timer fires at 2 min, this at 4 min
OnBootSec=4min
# Poll interval: every 5 min, same cadence as Suite
OnUnitActiveSec=5min
# Re-run missed firings after suspend/reboot
Persistent=true
Unit=drainq-one-mirror.service

[Install]
WantedBy=timers.target
