# ops/hetzner-update-proxy/drainq-one-mirror.service
[Unit]
Description=DrainQ ONE — Mirror releases from GitHub to local webroot
Documentation=file:///opt/drainq-mirror-one/DEPLOYMENT.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/drainq-mirror-one/mirror-releases-one.sh
User=www-data
Group=www-data

# Harden the service: no privilege escalation, filesystem largely read-only
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/www/drainq-updates/one /run
ReadOnlyPaths=/etc/drainq
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=drainq-one-mirror
